# com.bjorn.clash-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bjorn.clash-watchdog</string>
    <key>Program</key>
    <string>/Users/bjorn/ai-dev/clash-watchdog/clash-watchdog.sh</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/bjorn/.local/var/log/clash-watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bjorn/.local/var/log/clash-watchdog.log</string>
</dict>
</plist>
